# rootfs/etc/services.d/squid/run
#!/usr/bin/with-contenv sh

exec squid -N -f /config/etc/squid.conf
